# tasks/pingtest/pingtest.cron
#run task config crontab
#using run-parts to run prober and uploader sequently
*/2 * * * * run-parts /home/collect/tasks/pingtest
