# config/run.cron
*/2 * * * * run-parts /mbox/tasks/pingtest
